# fixtures/cron.d/every_minute
* * * * * randouser echo "every_minute"
